# packaging/debian/postinst
#!/bin/bash
# UCM Post-Installation Script - Simplified version
set -e

UCM_USER="ucm"
UCM_GROUP="ucm"
UCM_HOME="/opt/ucm"
UCM_DATA="/opt/ucm/data"
UCM_CONFIG="/etc/ucm"

case "$1" in
    configure)
        # Create user/group
        if ! getent group "$UCM_GROUP" >/dev/null 2>&1; then
            addgroup --system "$UCM_GROUP" >/dev/null
        fi
        if ! getent passwd "$UCM_USER" >/dev/null 2>&1; then
            adduser --system --home "$UCM_HOME" --no-create-home \
                    --ingroup "$UCM_GROUP" --disabled-password \
                    --shell /bin/false "$UCM_USER" >/dev/null
        fi
        
        # Create directories
        mkdir -p "$UCM_DATA"/{ca,certs,private,crl,scep,backups,sessions}
        mkdir -p "$UCM_CONFIG"
        
        # Set ownership
        chown -R "$UCM_USER":"$UCM_GROUP" "$UCM_HOME"
        chown -R "$UCM_USER":"$UCM_GROUP" "$UCM_CONFIG"
        chmod 700 "$UCM_DATA"/{ca,certs,private,backups}
        
        # Generate config if missing
        if [ ! -f "$UCM_CONFIG/ucm.env" ]; then
            SECRET_KEY=$(python3 -c 'import secrets; print(secrets.token_hex(32))')
            JWT_SECRET=$(python3 -c 'import secrets; print(secrets.token_hex(32))')
            ADMIN_PASSWORD=$(python3 -c 'import secrets; print(secrets.token_urlsafe(16))')
            FQDN=$(hostname -f 2>/dev/null || hostname)
            
            cat > "$UCM_CONFIG/ucm.env" << EOF
# UCM Configuration
BASE_DIR=$UCM_HOME
DATA_DIR=$UCM_DATA
FQDN=$FQDN
HTTPS_PORT=8443
SECRET_KEY=$SECRET_KEY
JWT_SECRET_KEY=$JWT_SECRET
DATABASE_PATH=$UCM_DATA/ucm.db
HTTPS_CERT_PATH=$UCM_DATA/https_cert.pem
HTTPS_KEY_PATH=$UCM_DATA/https_key.pem
INITIAL_ADMIN_USERNAME=admin
INITIAL_ADMIN_PASSWORD=$ADMIN_PASSWORD
INITIAL_ADMIN_EMAIL=admin@$FQDN
EOF
            chmod 600 "$UCM_CONFIG/ucm.env"
            chown "$UCM_USER":"$UCM_GROUP" "$UCM_CONFIG/ucm.env"
            
            echo ""
            echo "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
            echo " UCM INSTALLED!"
            echo " Admin: admin / $ADMIN_PASSWORD"
            echo " ⚠️  CHANGE PASSWORD AFTER FIRST LOGIN!"
            echo "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
        fi
        
        # Generate HTTPS cert if missing
        if [ ! -f "$UCM_DATA/https_cert.pem" ]; then
            echo "Generating HTTPS certificate..."
            FQDN=$(hostname -f 2>/dev/null || hostname)
            IP=$(hostname -I | awk '{print $1}')
            
            cat > /tmp/ucm-ssl.cnf << SSLEOF
[req]
distinguished_name = req_distinguished_name
x509_extensions = v3_req
prompt = no
[req_distinguished_name]
CN = $FQDN
[v3_req]
keyUsage = critical, digitalSignature, keyEncipherment
extendedKeyUsage = serverAuth
subjectAltName = DNS:$FQDN,DNS:localhost,IP:$IP,IP:127.0.0.1
basicConstraints = critical, CA:FALSE
SSLEOF
            
            openssl req -x509 -newkey rsa:4096 -sha256 -nodes \
                -keyout "$UCM_DATA/https_key.pem" \
                -out "$UCM_DATA/https_cert.pem" \
                -days 365 -config /tmp/ucm-ssl.cnf -extensions v3_req 2>/dev/null
            
            rm -f /tmp/ucm-ssl.cnf
            chmod 600 "$UCM_DATA/https_key.pem"
            chmod 644 "$UCM_DATA/https_cert.pem"
            chown "$UCM_USER":"$UCM_GROUP" "$UCM_DATA/https_key.pem" "$UCM_DATA/https_cert.pem"
            echo "✓ Certificate generated"
        fi
        
        # Create venv
        if [ ! -d "$UCM_HOME/venv" ]; then
            echo "Creating Python environment..."
            cd "$UCM_HOME"
            python3 -m venv venv
            venv/bin/pip install --quiet --upgrade pip
            venv/bin/pip install --quiet -r requirements.txt
            chown -R "$UCM_USER":"$UCM_GROUP" venv
            echo "✓ Python environment ready"
        fi
        
        # Initialize database
        if [ ! -f "$UCM_DATA/ucm.db" ]; then
            echo "Initializing database..."
            cd "$UCM_HOME/backend"
            set -a; [ -f "$UCM_CONFIG/ucm.env" ] && . "$UCM_CONFIG/ucm.env"; set +a
            sudo -u "$UCM_USER" ../venv/bin/python init_db.py
            echo "✓ Database initialized"
        fi
        
        # Clean up lock file
        rm -f "$UCM_DATA/.db_init.lock" 2>/dev/null || true
        
        systemctl daemon-reload >/dev/null 2>&1 || true
        systemctl enable ucm.service >/dev/null 2>&1 || true
        
        FQDN=$(hostname -f 2>/dev/null || hostname)
        echo ""
        echo " Access: https://$FQDN:8443"
        echo " Start:  systemctl start ucm"
        echo ""
        ;;
esac

exit 0
