# packaging/debian/postinst
#!/bin/bash
# UCM Post-Installation Script
# Handles configuration and service setup after package installation

set -e

# Source debconf library
. /usr/share/debconf/confmodule

# Package information
PACKAGE="ucm"
UCM_USER="ucm"
UCM_GROUP="ucm"
UCM_HOME="/opt/ucm"
UCM_DATA="/var/lib/ucm"
UCM_CONFIG="/etc/ucm"
UCM_LOG="/var/log/ucm"

case "$1" in
    configure)
        # Get configuration from debconf
        db_get ucm/fqdn
        FQDN="$RET"
        
        db_get ucm/https_port
        HTTPS_PORT="$RET"
        
        db_get ucm/admin_username
        ADMIN_USER="$RET"
        
        db_get ucm/admin_password
        ADMIN_PASS="$RET"
        
        db_get ucm/smtp_enable
        SMTP_ENABLED="$RET"
        
        # Create UCM user and group if they don't exist
        if ! getent group "$UCM_GROUP" >/dev/null 2>&1; then
            addgroup --system "$UCM_GROUP" >/dev/null
        fi
        
        if ! getent passwd "$UCM_USER" >/dev/null 2>&1; then
            adduser --system --home "$UCM_HOME" --no-create-home \
                    --ingroup "$UCM_GROUP" --disabled-password \
                    --shell /bin/false "$UCM_USER" >/dev/null
        fi
        
        # Create directories
        mkdir -p "$UCM_DATA"/{cas,certs,backups,logs,temp}
        mkdir -p "$UCM_CONFIG"
        mkdir -p "$UCM_LOG"
        
        # Set permissions
        chown -R "$UCM_USER":"$UCM_GROUP" "$UCM_HOME"
        chown -R "$UCM_USER":"$UCM_GROUP" "$UCM_DATA"
        chown -R "$UCM_USER":"$UCM_GROUP" "$UCM_CONFIG"
        chown -R "$UCM_USER":"$UCM_GROUP" "$UCM_LOG"
        
        chmod 750 "$UCM_DATA"
        chmod 700 "$UCM_DATA"/{cas,certs,backups}
        chmod 755 "$UCM_LOG"
        
        # Generate secrets
        SECRET_KEY=$(python3 -c 'import secrets; print(secrets.token_hex(32))')
        JWT_SECRET=$(python3 -c 'import secrets; print(secrets.token_hex(32))')
        
        # Create .env configuration file
        cat > "$UCM_CONFIG/ucm.env" <<EOF
# UCM Configuration
# Generated by Debian package installer on $(date)

# Core Settings
FQDN=$FQDN
HTTPS_PORT=$HTTPS_PORT
HTTP_PORT=8080
DEBUG=false
LOG_LEVEL=INFO

# Security
SECRET_KEY=$SECRET_KEY
JWT_SECRET_KEY=$JWT_SECRET
SESSION_TIMEOUT=3600
JWT_EXPIRATION=86400

# Database
DATABASE_PATH=$UCM_DATA/ucm.db
BACKUP_ENABLED=true
BACKUP_RETENTION_DAYS=30

# HTTPS Certificate
HTTPS_CERT_PATH=$UCM_DATA/https_cert.pem
HTTPS_KEY_PATH=$UCM_DATA/https_key.pem
HTTPS_AUTO_GENERATE=true

# Initial Admin User
INITIAL_ADMIN_USERNAME=$ADMIN_USER
INITIAL_ADMIN_PASSWORD=$ADMIN_PASS
INITIAL_ADMIN_EMAIL=admin@$FQDN

# Caching
CACHE_ENABLED=true
CACHE_TYPE=simple
CACHE_DEFAULT_TIMEOUT=300

# Protocols
ACME_ENABLED=true
SCEP_ENABLED=true

# Logging
LOG_FILE=$UCM_LOG/ucm.log
AUDIT_LOG_FILE=$UCM_LOG/audit.log
EOF

        # Add SMTP configuration if enabled
        if [ "$SMTP_ENABLED" = "true" ]; then
            db_get ucm/smtp_server
            SMTP_SERVER="$RET"
            db_get ucm/smtp_port
            SMTP_PORT="$RET"
            db_get ucm/smtp_username
            SMTP_USER="$RET"
            db_get ucm/smtp_password
            SMTP_PASS="$RET"
            db_get ucm/smtp_from
            SMTP_FROM="$RET"
            
            cat >> "$UCM_CONFIG/ucm.env" <<EOF

# SMTP Configuration
SMTP_ENABLED=true
SMTP_SERVER=$SMTP_SERVER
SMTP_PORT=$SMTP_PORT
SMTP_USERNAME=$SMTP_USER
SMTP_PASSWORD=$SMTP_PASS
SMTP_FROM=$SMTP_FROM
SMTP_TLS=true
EOF
        else
            cat >> "$UCM_CONFIG/ucm.env" <<EOF

# SMTP Configuration
SMTP_ENABLED=false
EOF
        fi
        
        # Secure the configuration file
        chmod 600 "$UCM_CONFIG/ucm.env"
        chown "$UCM_USER":"$UCM_GROUP" "$UCM_CONFIG/ucm.env"
        
        # Create symlink from /opt/ucm to config
        ln -sf "$UCM_CONFIG/ucm.env" "$UCM_HOME/.env" 2>/dev/null || true
        
        # Generate HTTPS certificate if requested
        db_get ucm/generate_https_cert
        if [ "$RET" = "true" ]; then
            if [ ! -f "$UCM_DATA/https_cert.pem" ]; then
                echo "Generating self-signed HTTPS certificate..."
                
                # Create OpenSSL config with SAN
                cat > /tmp/ucm-ssl.cnf <<SSLEOF
[req]
distinguished_name = req_distinguished_name
x509_extensions = v3_req
prompt = no

[req_distinguished_name]
C = US
ST = State
L = City
O = Ultimate CA Manager
OU = IT
CN = $FQDN

[v3_req]
keyUsage = keyEncipherment, dataEncipherment
extendedKeyUsage = serverAuth
subjectAltName = @alt_names

[alt_names]
DNS.1 = $FQDN
DNS.2 = localhost
IP.1 = 127.0.0.1
SSLEOF
                
                openssl req -x509 -newkey rsa:4096 -nodes \
                    -keyout "$UCM_DATA/https_key.pem" \
                    -out "$UCM_DATA/https_cert.pem" \
                    -days 365 \
                    -config /tmp/ucm-ssl.cnf \
                    -extensions v3_req \
                    2>/dev/null
                
                rm /tmp/ucm-ssl.cnf
                
                chmod 600 "$UCM_DATA/https_key.pem"
                chmod 644 "$UCM_DATA/https_cert.pem"
                chown "$UCM_USER":"$UCM_GROUP" "$UCM_DATA/https_key.pem"
                chown "$UCM_USER":"$UCM_GROUP" "$UCM_DATA/https_cert.pem"
                
                echo "âœ“ HTTPS certificate generated"
            fi
        fi
        
        # Reload systemd
        systemctl daemon-reload >/dev/null 2>&1 || true
        
        # Enable and start service if requested
        db_get ucm/autostart
        if [ "$RET" = "true" ]; then
            systemctl enable ucm.service >/dev/null 2>&1 || true
            
            if systemctl is-active --quiet ucm.service; then
                systemctl restart ucm.service >/dev/null 2>&1 || true
            else
                systemctl start ucm.service >/dev/null 2>&1 || true
            fi
            
            # Wait for service to start
            sleep 2
            
            if systemctl is-active --quiet ucm.service; then
                echo "âœ“ UCM service started successfully"
            else
                echo "âš  UCM service failed to start. Check logs: journalctl -u ucm"
            fi
        fi
        
        # Display completion message
        db_subst ucm/configuration_complete fqdn "$FQDN"
        db_subst ucm/configuration_complete port "$HTTPS_PORT"
        db_subst ucm/configuration_complete admin_username "$ADMIN_USER"
        db_input high ucm/configuration_complete || true
        db_go || true
        
        # Clean up debconf database (keep FQDN for reconfiguration)
        db_reset ucm/admin_password
        db_reset ucm/smtp_password
        
        echo ""
        echo "â•”â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•—"
        echo "â•‘   UCM Installation Complete! ðŸŽ‰       â•‘"
        echo "â•šâ•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•"
        echo ""
        echo "Access: https://$FQDN:$HTTPS_PORT"
        echo "Login:  $ADMIN_USER / (your password)"
        echo ""
        echo "Configuration: $UCM_CONFIG/ucm.env"
        echo "Logs:          journalctl -u ucm -f"
        echo "Service:       systemctl status ucm"
        echo ""
        echo "To reconfigure: dpkg-reconfigure ucm"
        echo ""
        
        ;;
        
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
        
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# Stop debconf
db_stop

exit 0
