# packaging/debian/ucm.service
[Unit]
Description=Ultimate CA Manager (UCM)
Documentation=https://github.com/NeySlim/ultimate-ca-manager
After=network.target

[Service]
Type=simple
User=ucm
Group=ucm
WorkingDirectory=/opt/ucm/backend

# Environment
EnvironmentFile=/etc/ucm/ucm.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536

# Execution
ExecStart=/opt/ucm/venv/bin/gunicorn \
    --config /opt/ucm/gunicorn.conf.py \
    --bind 0.0.0.0:${HTTPS_PORT} \
    --certfile /opt/ucm/data/https_cert.pem \
    --keyfile /opt/ucm/data/https_key.pem \
    --access-logfile /opt/ucm/data/access.log \
    --error-logfile /opt/ucm/data/error.log \
    wsgi:app
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
