# debian-test/ucm-1.8.4/DEBIAN/prerm
#!/bin/bash
set -e
systemctl stop ucm.service || true
systemctl disable ucm.service || true
exit 0
